# freebsd/webapp/supervisord.conf
; freebsd/webapp/supervisord.conf
;
; Base (non-station) supervisord config for the `webapp` jail. Install to:
;   /usr/local/etc/supervisord.conf
;
; Translated from util/docker/supervisor/supervisor/supervisord.conf +
; util/docker/web/service.full/{nginx,php-fpm,centrifugo,sftpgo}.conf.
;
; Deliberately NOT included here (out of scope for this change):
;   - [program:cron] — cron is handled by FreeBSD's native cron(8) +
;     freebsd/webapp/crontab instead of supercronic, so there is no
;     supervisord program for it.
;   - [program:php-worker] / [program:php-nowplaying] — these run
;     `azuracast_cli queue:process` / `azuracast_cli azuracast:sync:nowplaying`
;     against the deployed PHP application, whose install path/name this
;     change does not define (app code deployment is out of scope). Add
;     them here once that path is known, following the nginx/php-fpm
;     stanzas below as a template (autostart=true, priority=600,
;     user=azuracast, stopasgroup/killasgroup=true).
;   - Per-station backend/frontend programs — generated dynamically by
;     the application itself into two separate files per station,
;     /var/azuracast/stations/*/config/supervisord.{backend,frontend}.conf
;     (see [include] below; split into two files rather than one
;     combined supervisord.conf so a station's frontend/backend can each
;     be reloaded/managed independently, including when frontend runs in
;     a separate jail from backend — see Configuration.php's
;     buildAdapterSupervisorConfig()); not managed by this file.
;   - [program:redis]/valkey — NOT managed by supervisord here. Valkey
;     ships its own FreeBSD rc.d script (`pkg install valkey` registers
;     an rc.d service); enable it independently with
;     `sysrc valkey_enable=YES` and let rc(8) manage it, rather than
;     duplicating that under supervisord as Docker's minimal.conf.d does.
;
; Startup ORDERING note: Docker's app_startup script explicitly sequences
; php-fpm+centrifugo -> (sleep 2s) -> nginx -> "start all" because
; supervisord starts life with everything at autostart=false and a
; dedicated "startup" program drives the sequence via supervisorctl.
; On FreeBSD, freebsd/webapp/rc.d/azuracast already waits for MariaDB and
; runs one-time setup/migrations *before* supervisord is even started, so
; there's no equivalent race to guard against here — every program below
; is autostart=true and supervisord's own ascending-priority start order
; (php-fpm 300 < nginx 400 < centrifugo/sftpgo 700) is sufficient. nginx
; proxies to centrifugo over TCP and will simply retry/502 for the first
; moment if centrifugo isn't up yet, then recover — same tolerance the
; Docker setup relies on implicitly.

; The socket path is NOT a free choice: backend/config/services.php
; hardcodes /var/run/supervisor.sock (CURLOPT_UNIX_SOCKET_PATH) for the
; local supervisor XML-RPC client -- moving it (as an earlier workaround
; did, when the rc.subr collision documented in rc.d/azuracast wrongly
; ran supervisord as non-root and /var/run was unwritable) silently
; breaks every station operation with "cURL error 7 ... unix:///var/run/
; supervisor.sock" in the app log. supervisord runs as root now, so the
; upstream path just works; chmod/chown below keep it reachable by the
; azuracast-uid php-fpm workers that actually connect to it.
[unix_http_server]
file = /var/run/supervisor.sock
chmod = 0666
chown = azuracast:azuracast

[supervisord]
; supervisord itself now genuinely starts as root under rc.d (the
; freebsd/webapp/rc.d/azuracast collision that silently ran it as
; azuracast instead -- see that file's azuracast_setup_user comment --
; is fixed), so this "drop privilege to root" (a no-op, since it's
; already root) is safe again. It's what lets nginx bind ports <1024
; (nginx/php-fpm then privilege-drop themselves internally, via
; nginx.conf's "user azuracast;" and php-fpm's own pool config) and lets
; the explicit "user = azuracast" directives on the centrifugo/sftpgo
; programs below actually take effect (setuid down from root), instead
; of failing with "Can't drop privilege as nonroot user".
user = root
logfile = /var/log/azuracast/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /var/run/supervisord.pid
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///var/run/supervisor.sock

; Pulls in per-station configs generated by the application at runtime.
; Deliberately matches ONLY supervisord.backend.conf -- the frontend
; (Icecast) runs in a separate per-station jail in this deployment, and
; its supervisord.frontend.conf is consumed by THAT jail's supervisord
; (freebsd/icecast/supervisord.conf.tmpl) via the shared nullfs mount.
; A wildcard supervisord.*.conf here (as an earlier revision had) makes
; THIS jail's supervisord also load the frontend program and try to exec
; an Icecast binary that isn't installed here -- confirmed on a real
; install the moment the first station's configs were written.
[include]
files = /usr/local/etc/supervisor.d/*.conf /var/azuracast/stations/*/config/supervisord.backend.conf

; --- nginx --------------------------------------------------------------
; Docker: util/docker/web/service.full/nginx.conf
[program:nginx]
command = /usr/local/sbin/nginx -g "daemon off;" -c /usr/local/etc/nginx/nginx.conf
priority = 400
numprocs = 1
autostart = true
autorestart = true

stopasgroup = true
killasgroup = true

stdout_logfile = /var/log/azuracast/nginx.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5
redirect_stderr = true

; --- php-fpm --------------------------------------------------------------
; Docker: util/docker/web/service.full/php-fpm.conf
; FreeBSD's php85 port installs php-fpm to /usr/local/sbin/php-fpm and its
; default pool config to /usr/local/etc/php-fpm.conf +
; /usr/local/etc/php-fpm.d/*.conf (equivalent of Debian's php-fpm.d).
[program:php-fpm]
command = /usr/local/sbin/php-fpm -F --fpm-config /usr/local/etc/php-fpm.conf
priority = 300
numprocs = 1
autostart = true
autorestart = true

stopasgroup = true
killasgroup = true

stdout_logfile = /var/log/azuracast/php-fpm.log
stdout_logfile_maxbytes = 0
stderr_logfile = /var/log/azuracast/php-fpm.error.log
stderr_logfile_maxbytes = 0

; --- centrifugo --------------------------------------------------------------
; Docker: util/docker/web/service.full/centrifugo.conf
[program:centrifugo]
directory = /var/azuracast/centrifugo
command = /usr/local/bin/centrifugo -c /var/azuracast/centrifugo/config.toml
user = azuracast
priority = 700
numprocs = 1
autostart = true
autorestart = true

stopasgroup = true
killasgroup = true

stdout_logfile = /var/log/azuracast/centrifugo.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5
redirect_stderr = true

; --- sftpgo --------------------------------------------------------------
; Docker: util/docker/web/service.full/sftpgo.conf
[program:sftpgo]
directory = /var/azuracast/sftpgo
command = /usr/local/bin/sftpgo --config-dir=/var/azuracast/sftpgo serve -l ""
user = azuracast
priority = 700
numprocs = 1
autostart = true
autorestart = true

stopasgroup = true
killasgroup = true

stdout_logfile = /var/log/azuracast/sftpgo.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5
redirect_stderr = true
